# client/iprep.cron
# Start iprep
*/5 *   * * *  root    /usr/local/bin/iprep/start_iprep.sh
