# client/iprep.cron
# Start iprep
*/5 *   * * *   zimbra    /usr/local/bin/iprep/start_iprep.sh
