# server/iprep.cron
*/1 *   * * *   root /usr/local/iprep/bin/iprep.pl; /usr/bin/pkill -HUP named
